# supervisord-dump1090.ini
[program:dump1090]
command=/usr/local/bin/dump1090 --net-only --write-json /usr/share/nginx/html/data
autostart=true
autorestart=false
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
